# srcs/requirements/bonus/redis/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:redis-server]
command=redis-server /etc/redis/redis.conf